# debian/usbhasp.usbhaspd.service
[Unit]
Description=Usbhasp daemon

[Service]
Type=simple
ExecStart=/usr/sbin/usbhaspd.sh

[Install]
WantedBy=multi-user.target
 
